# build/bin/install
#!/bin/sh
export PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin

# ===================================================================

cat << EOF >> /etc/resolv.conf
nameserver 192.168.1.1
nameserver 1.1.1.1
nameserver 8.8.8.8
EOF

# ===================================================================

apk add vim openssh dhcpcd bash openrc xterm podman

# ===================================================================

cat << 'EOF' > /etc/init.d/shell
#!/sbin/openrc-run

command=/bin/sh /bin/init

start() {
    ebegin "Starting shell"
    exec ${command}
    eend $?
}
EOF
chmod +x /etc/init.d/shell

# ===================================================================

rc-update add sshd default
rc-update add dhcpcd default
rc-update add shell default

# ===================================================================

cat << EOF > /root/.bashrc
PS1='${debian_chroot:+($debian_chroot)}\u@\h:\w# '
EOF

# ===================================================================

cat << EOF > /etc/fstab
proc            /proc           proc    defaults        0       0
tmpfs           /run            tmpfs   mode=0755,nosuid,nodev  0       0
sysfs           /sys            sysfs   defaults        0       0
EOF

# ===================================================================

mkdir -p /run /run/lock /var/run
